# files/general/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
exec php-fpm --nodaemonize
